# etc/cron.d/udge
# /etc/cron.d/udge
#
# Udge - GPLv2 - Copyright (C) 2020 Rudy Matela

# Check run submissions, pick new ones, update user HTMLs.
* * * * * udge /usr/local/bin/udge-check-and-pick 15 15

# Run submissions.
* * * * * udge-1 /usr/local/bin/udge-compile-and-run 1
* * * * * udge-2 /usr/local/bin/udge-compile-and-run 2
* * * * * udge-3 /usr/local/bin/udge-compile-and-run 3
* * * * * udge-4 /usr/local/bin/udge-compile-and-run 4
* * * * * udge-5 /usr/local/bin/udge-compile-and-run 5
* * * * * udge-6 /usr/local/bin/udge-compile-and-run 6

# Update the rank HTML.
*/2 * * * * udge /usr/local/bin/udge-update-rank-html

# Make a backup every day at 3:33.
33 3 * * * udge /usr/local/bin/udge-backup
